# config/rpi-w1-ds2482/ds2482-probe.service
[Unit]
Description=Probe DS248x I2C-1wire interface device

[Service]
ExecStart=/sbin/ds2482-probe

[Install]
WantedBy=default.target
